# highscores-api/config/systemd/highscores-api.service
[Unit]
Description=Highscores API Service
After=network.target nginx.service
Requires=network.target
Wants=nginx.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app/Project-2025/highscores-api
ExecStart=/home/ubuntu/app/Project-2025/highscores-api/scripts/start_server.sh
ExecStartPre=/bin/mkdir -p /var/log/highscores-api
ExecStartPre=/bin/chown -R ubuntu:ubuntu /var/log/highscores-api
ExecStartPre=/bin/mkdir -p /home/ubuntu/app/Project-2025/highscores-api/storage
ExecStartPre=/bin/chmod -R 755 /home/ubuntu/app/Project-2025/highscores-api/storage

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Environment variables
Environment=RAILS_ENV=production
Environment=RAILS_SERVE_STATIC_FILES=true
Environment=RAILS_LOG_TO_STDOUT=true

# Logging
StandardOutput=append:/var/log/highscores-api/application.log
StandardError=append:/var/log/highscores-api/error.log

# Resource limits
TimeoutStartSec=300
TimeoutStopSec=90
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
